# crontab
# START CRON JOB
0 5 * * * python3 /app/soliscloud_to_pvoutput.py > /proc/1/fd/1 2>/proc/1/fd/2
# END CRON JOB